# cookbooks/web_api/templates/defaults/flask.service.erb
[Unit]
Description=Gunicorn instance to serve Flask app
After=network.target

[Service]
User=root
WorkingDirectory=/opt/flask_app
ExecStart=/usr/bin/gunicorn --workers 3 --bind 127.0.0.1:8000 app:app

[Install]
WantedBy=multi-user.target
